# luci-app-modem/root/etc/init.d/modeminit
#!/bin/sh /etc/rc.common

START=99

start() {
  # /bin/sh /usr/share/cpe/rssi &
  /bin/sh /usr/share/modem/modem_data &
}